# templates/ipfs-gc.timer.j2
[Timer]
OnCalendar={{ ipfs_gc_when }}

[Install]
WantedBy=timers.target
